# infrastructure/automation/job-daily-run.service
[Unit]
Description=Daily Job Application Automation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=%i
WorkingDirectory=/home/jclee/dev/resume/apps/job-server
ExecStart=/usr/local/go/bin/go run ./scripts/ops/auto-daily-run/main.go
Environment=PATH=/usr/local/go/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/%i
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
